# BlogCloudServer/bin/blog.service
# put this file in /usr/lib/systemd/system/blog.service as root
# below paths assume bedock installed in /opt/blog
# systemctl start blog


[Unit]
Description=Blog Cloud
After=network.target

[Service]
Type=forking
Environment="JAVA_HOME=/usr/local/tools/java_jdk"
ExecStart=/opt/Blog/bin/run.sh
PrivateTmp=true
#Restart=always

[Install]
WantedBy=multi-user.target
